# services/librechat/launchd/com.ai.librechat.mongodb.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ai.librechat.mongodb</string>

    <key>UserName</key>
    <string>librechat</string>

    <key>WorkingDirectory</key>
    <string>/var/lib/librechat</string>

    <key>ProgramArguments</key>
    <array>
      <string>/var/lib/librechat/bin/mongod</string>
      <string>--config</string>
      <string>/var/lib/librechat/mongo/mongod.conf</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/var/lib/librechat</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/librechat/mongodb.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/librechat/mongodb.err.log</string>
  </dict>
</plist>
