# services/followyourcanvas/launchd/com.followyourcanvas.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.followyourcanvas</string>

    <key>UserName</key><string>followyourcanvas</string>
    <key>WorkingDirectory</key><string>/var/lib/followyourcanvas/app</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>source /var/lib/followyourcanvas/followyourcanvas.env 2&gt;/dev/null || true; exec ${FYC_CMD}</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key><string>/var/lib/followyourcanvas</string>
      <key>PYTHONUNBUFFERED</key><string>1</string>
    </dict>

    <key>StandardOutPath</key><string>/var/log/followyourcanvas/followyourcanvas.out.log</string>
    <key>StandardErrorPath</key><string>/var/log/followyourcanvas/followyourcanvas.err.log</string>
  </dict>
</plist>
