# services/pocket-tts/launchd/com.pocket-tts.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.pocket-tts.server</string>

    <key>UserName</key><string>pockettts</string>
    <key>WorkingDirectory</key><string>/var/lib/pocket-tts</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>source /etc/pocket-tts/pocket-tts.env 2&gt;/dev/null || true; exec /var/lib/pocket-tts/env/bin/uvicorn pocket_tts_server:app --host "${POCKET_TTS_HOST:-0.0.0.0}" --port "${POCKET_TTS_PORT:-9940}"</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key><string>/var/lib/pocket-tts</string>
      <key>PYTHONUNBUFFERED</key><string>1</string>
    </dict>

    <key>StandardOutPath</key><string>/var/log/pocket-tts/pocket-tts.out.log</string>
    <key>StandardErrorPath</key><string>/var/log/pocket-tts/pocket-tts.err.log</string>
  </dict>
</plist>
