# services/skyreels-v2/launchd/com.skyreels-v2.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.skyreels-v2.server</string>

    <key>UserName</key><string>skyreels</string>
    <key>WorkingDirectory</key><string>/var/lib/skyreels-v2</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>source /etc/skyreels-v2/skyreels-v2.env 2&gt;/dev/null || true; exec /var/lib/skyreels-v2/venv/bin/uvicorn skyreels_shim_server:app --host "${SKYREELS_HOST:-127.0.0.1}" --port "${SKYREELS_PORT:-9190}"</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key><string>/var/lib/skyreels-v2</string>
      <key>PYTHONUNBUFFERED</key><string>1</string>
    </dict>

    <key>StandardOutPath</key><string>/var/log/skyreels-v2/skyreels-v2.out.log</string>
    <key>StandardErrorPath</key><string>/var/log/skyreels-v2/skyreels-v2.err.log</string>
  </dict>
</plist>
