# services/heartmula/launchd/com.heartmula.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.heartmula.server</string>

    <key>UserName</key>
    <string>heartmula</string>

    <key>ProgramArguments</key>
    <array>
      <string>/var/lib/heartmula/env/bin/python</string>
      <string>/var/lib/heartmula/heartmula_server.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HEARTMULA_MODEL_PATH</key>
      <string>/var/lib/heartmula/ckpt</string>
      <key>HEARTMULA_OUTPUT_DIR</key>
      <string>/var/lib/heartmula/output</string>
      <key>HEARTMULA_HOST</key>
      <string>127.0.0.1</string>
      <key>HEARTMULA_PORT</key>
      <string>9920</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/var/lib/heartmula</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/heartmula/heartmula.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/heartmula/heartmula.err.log</string>
  </dict>
</plist>
